# systemd.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
ExecStart=/bin/node /home/roboshop/cart/server.js

[Install]
WantedBy=multi-user.target